# deploy/robotunneld.service
[Unit]
Description=RoboTunnel Daemon (robotunneld)
Documentation=https://github.com/get-robotunnel/robot-agent-tunnel
After=network.target
Wants=network.target

[Service]
Type=simple
User=robotunnel
Group=robotunnel
ExecStart=/usr/local/bin/robotunneld
Restart=on-failure
RestartSec=5s

# IPC socket directory
RuntimeDirectory=robotunnel
RuntimeDirectoryMode=0750

# Environment
EnvironmentFile=-/etc/robotunnel/daemon.env
Environment=RT_DAEMON_SOCKET=/run/robotunnel/rt.sock
Environment=RT_DAEMON_LISTEN_PORT=11411
Environment=RT_DAEMON_INSECURE=false

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/robotunnel

[Install]
WantedBy=multi-user.target
