# deploy/tunnel.service
[Unit]
Description=RoboTunnel Tunnel Service (connection layer)
After=network.target

[Service]
Type=simple
User=robotunnel
WorkingDirectory=/opt/robotunnel-tunnel
ExecStart=/opt/robotunnel-tunnel/bin/tunnel-svc
Restart=always
RestartSec=5
EnvironmentFile=/opt/robotunnel-tunnel/config/.env
StandardOutput=append:/var/log/robotunnel-tunnel.log
StandardError=append:/var/log/robotunnel-tunnel.log

[Install]
WantedBy=multi-user.target
